# deploy/labsoldier.service
[Unit]
Description=LabSoldier backend (Express + SQLite)
After=network.target

[Service]
Type=simple
# アプリの配置先（setup.sh で /opt/labsoldier にcloneする前提）
WorkingDirectory=/opt/labsoldier/backend
# Node のフルパスは `which node` で確認して合わせる
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=3
User=labsoldier
Group=labsoldier
Environment=NODE_ENV=production
Environment=PORT=3001
# 永続ディスク上の DB パス
Environment=DB_PATH=/opt/labsoldier/data/labsoldier.db
# 環境変数は EnvironmentFile にまとめる（JWT_SECRET / LAB_ALLOWED_IPS など）
EnvironmentFile=/opt/labsoldier/labsoldier.env

[Install]
WantedBy=multi-user.target
